# configs/etc/systemd/system/scw-generate-ssh-keys.service.d/sckl-disable.conf
[Unit]
ConditionKernelCommandLine=!sckl-kernel-loaded
